# roles/prometheus/templates/prometheus.service
[Unit]
Description=Prometheus Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/prometheus-2.53.0.linux-amd64/prometheus
WorkingDirectory=/opt

[Install]
WantedBy=multi-user.target